# cron/crontab
25 6	* * *	root	cd / && run-parts --report /var/www/kreddy.ru.new/protected/cron/sh/